# catfinder.service
[Unit]
Description=CatFinder PRO RTSP cat detector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/catfinder_pro
EnvironmentFile=/opt/catfinder_pro/.env
ExecStart=/opt/catfinder_pro/.venv/bin/python -m app.main
Restart=always
RestartSec=5
User=catfinder
Group=catfinder

[Install]
WantedBy=multi-user.target
